# deploy/powertracker.service
[Unit]
Description=Power Tracker Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=colbyharris
WorkingDirectory=/home/colbyharris/Documents/PowerTrackerService
EnvironmentFile=/home/colbyharris/Documents/PowerTrackerService/.env
ExecStart=/usr/bin/python3 -m power_tracker.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
